# debian/bhfuse.upstart
description "BitHorde FUSE proxy"

start on started bithorde
stop on stopped bithorde

pre-start script
  if [ ! -d /tmp/bhfuse ]; then
    mkdir /tmp/bhfuse
  fi
end script

exec /usr/bin/bhfuse --mountpoint /tmp/bhfuse
respawn

post-stop script
 (mount | grep -c /tmp/bhfuse) && /bin/umount -lf /tmp/bhfuse
 sleep 1
end script
