# debian/bithorded.bithorded-tcp.socket
[Socket]
ListenStream=1337
FileDescriptorName=tcp
Service=bithorded.service

[Install]
WantedBy=sockets.target
