# debian/bithorde.upstart


start on net-device-up


exec /usr/bin/bithorded /etc/bithorde.conf

expect fork


