# debian/bithorde.upstart

start on (filesystem and net-device-up IFACE!=lo)

setuid bithorde
console log

exec /usr/bin/bithorded -c/etc/bithorde.conf
respawn
respawn limit 5 60

