# debian/bithorde.service
[Unit]
Description=BitHorde daemon

[Service]
ExecStart=/usr/bin/bithorded -c/etc/bithorde.conf

[Install]
WantedBy=multi-user.target
