# debian/bithorded.bithorded-unix.socket
[Socket]
ListenStream=/var/run/bithorde/socket
FileDescriptorName=unix
Service=bithorded.service

[Install]
WantedBy=sockets.target
